# supervisord-remoting.conf
[program:remoting-server]
command=java -jar /remoting.jar
numprocs=1
autostart=true
autorestart=false
startsecs=5